# conf.d/bundler.conf
[program:bundler-install]
directory=repos/Codecademy
command=%(ENV_HOME)s/.rbenv/shims/bundle install --binstubs --path=.bundle
autostart=true
autorestart=false
