# conf.d/codex.conf
[program:codex]
directory=repos/codex-adapter
command=bash script/bootstrap

