# conf.d/contenteditor-gulp.conf
[program:contenteditor-gulp]
directory=repos/content-editor/
command=%(ENV_NVM_BIN)s/node node_modules/gulp/bin/gulp.js
autostart=false
