# conf.d/content_service.conf
[program:content_service]
directory=repos/content_service
command=bundle exec rackup
autostart=true
