# conf.d/contenteditor-app.conf
[program:content-editor-app]
directory=repos/content-editor
command=./start.sh
environment=CC_CONTENT_PATH='../Codecademy-content'
autostart=true
autorestart=false
redirect_stderr=true
