# conf.d/elasticsearch.conf
[program:elasticsearch]
command=/usr/local/bin/elasticsearch -f
autostart=true
autorestart=true
