# conf.d/mailcatcher.conf
[program:mailcatcher]
directory=repos/Codecademy
command=%(ENV_HOME)s/.rbenv/shims/mailcatcher
autostart=true
autorestart=false
