# conf.d/memcached.conf
[program:memcached]
command=/usr/local/bin/memcached
autostart=true
autorestart=true

