# conf.d/mongo.conf
[program:mongo]
command=/usr/local/bin/mongod
autostart=true
autorestart=true
