# conf.d/npm.conf
[program:npm-install]
directory=repos/Codecademy
command=%(ENV_NVM_BIN)s/npm install
autostart=true
autorestart=false
