# conf.d/programseditor-app.conf
[program:programs-editor-app]
directory=repos/programs-editor
command=%(ENV_NVM_BIN)s/npm run start
autostart=true
autorestart=false
redirect_stderr=true
