# conf.d/propeller.conf
[program:propeller]
directory=%(ENV_GOPATH)s/src/github.com/RyzacInc/workshop/
command=/bin/bash -c "docker-compose down && docker-compose up"
autostart=true
autorestart=true
