# conf.d/rails.conf
[program:rails]
directory=repos/Codecademy
command=%(ENV_HOME)s/.rbenv/shims/bundle exec rails server
environment=RAILS_ENV='development'
autostart=true
autorestart=false
