# conf.d/webpack.conf
[program:webpack]
directory=repos/Codecademy
command=%(ENV_NVM_BIN)s/npm start
